# Dark Mode/com.user.darkmode.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <!-- A unique identifier for this job -->
    <key>Label</key>
    <string>com.user.darkmode</string>
    
    <!-- Run the script every 60 seconds -->
    <key>StartInterval</key>
    <integer>60</integer>
    
    <!-- Specify the full path to your script -->
    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/darkmode.sh</string>
    </array>
    
    <!-- Optional: Standard output and error paths for logging -->
    <key>StandardOutPath</key>
    <string>/tmp/darkmode.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/darkmode.err</string>
    
    <!-- Run the job immediately on load -->
    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>